# distro-server/scripts/amplifier-distro.service
[Unit]
Description=Amplifier Distro Server
After=network.target

[Service]
Type=simple
ExecStart=amp-distro serve --port 8400
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60
WorkingDirectory=%h
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
